# debian/prerm
#!/bin/sh
set -e
if [ "$1" = remove ] ; then
	deb-systemd-helper disable novena-disable-ssp.service
fi
#DEBHELPER#
exit 0
